# services/ai.meta-factory.sage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Sage daemon launchd unit.

  Installed by arc 0.26.0's darwin-launchd HostAdapter into
  ~/Library/LaunchAgents/ai.meta-factory.sage.plist. Arc renders the
  following tokens at install time (see launchd-install.ts buildLaunchdTokens):

    {{BIN}}           → absolute path of ~/bin/sage (the launcher symlink)
    {{INSTALL_PATH}}  → the cloned package directory
    {{HOME}}          → os.homedir()
    {{LOG_DIR}}       → ~/Library/Logs/sage/
    {{NATS_URL}}      → process.env.NATS_URL or nats://127.0.0.1:4222

  Sage-specific config (org, did, source) is hardcoded here — only the
  five well-known tokens above are substituted by arc.

  Provider API keys are deliberately NOT in this plist. They live in the
  operator's login shell env (or a per-user keychain) and reach the pi
  subprocess via Sage's allow-list (src/pi/env.ts). Embedding keys in a
  world-readable plist would be a security regression.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.meta-factory.sage</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{BIN}}</string>
    <string>serve</string>
    <string>--nats</string>
    <string>{{NATS_URL}}</string>
    <string>--org</string>
    <string>metafactory</string>
    <string>--did</string>
    <string>did:mf:sage</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>NATS_URL</key>
    <string>{{NATS_URL}}</string>
    <key>NATS_CREDS_FILE</key>
    <string>{{HOME}}/.config/nats/creds/sage.creds</string>
    <key>SAGE_ORG</key>
    <string>metafactory</string>
    <key>SAGE_DID</key>
    <string>did:mf:sage</string>
    <key>SAGE_SOURCE</key>
    <string>metafactory.sage.local</string>
    <key>SAGE_DATA_RESIDENCY</key>
    <string>CH</string>
    <key>PI_PROVIDER</key>
    <string>openrouter</string>
  </dict>

  <key>StandardOutPath</key>
  <string>{{LOG_DIR}}/sage.log</string>
  <key>StandardErrorPath</key>
  <string>{{LOG_DIR}}/sage.err.log</string>

  <key>WorkingDirectory</key>
  <string>{{INSTALL_PATH}}</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
